# config/ld.so.preload
/usr/local/lib/libhardened_malloc-light.so
